# tailscale/usr/lib/systemd/system/tailscaled.service.d/50-sysext-default.conf
# Set defaults from /usr/etc/default/tailscaled
[Service]
EnvironmentFile=
EnvironmentFile=/usr/etc/default/tailscaled
