# nordvpn-gui/usr/lib/systemd/system/nordvpnd.service.d/librarypath.conf
[Service]
Environment=LD_LIBRARY_PATH="/usr/lib/nordvpn/"
