# qemu-guest-agent/usr/lib/systemd/system/qemu-guest-agent.service.d/override.conf
[Service]
EnvironmentFile=
EnvironmentFile=/usr/etc/sysconfig/qemu-ga
EnvironmentFile=-/etc/sysconfig/qemu-ga
